# firefly-iii/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Firefly III
# Runs the PHP-FPM daemon
# ==============================================================================

declare key

key=$(cat /data/firefly/appkey.txt)
export APP_KEY=${key}
export LOG_CHANNEL=stdout
export APP_ENV=local
export DB_CONNECTION=mysql
export DB_DATABASE
export DB_HOST
export DB_PASSWORD
export DB_PORT
export DB_USERNAME
export MAPBOX_API_KEY

if bashio::config.has_value 'mapbox_api_key';then
    MAPBOX_API_KEY=$(bashio::config "mapbox_api_key")
fi

if bashio::config.has_value 'remote_mysql_host';then
    DB_HOST=$(bashio::config "remote_mysql_host")
    DB_DATABASE=$(bashio::config "remote_mysql_database")
    DB_USERNAME=$(bashio::config "remote_mysql_username")
    DB_PASSWORD=$(bashio::config "remote_mysql_password")
    DB_PORT=$(bashio::config "remote_mysql_port")
else
    DB_HOST=$(bashio::services "mysql" "host")
    DB_DATABASE=firefly
    DB_USERNAME=$(bashio::services "mysql" "username")
    DB_PASSWORD=$(bashio::services "mysql" "password")
    DB_PORT=$(bashio::services "mysql" "port")
fi

if bashio::config.equals 'log_level' 'debug' \
    || bashio::config.equals 'log_level' 'trace' ;then
    export APP_DEBUG=true
    export APP_LOG_LEVEL=debug
fi
bashio::log.info "Running migration commands..."
php /var/www/firefly/artisan firefly-iii:create-database
php /var/www/firefly/artisan migrate --seed --no-interaction --force
#php /var/www/firefly/artisan firefly-iii:decrypt-all

bashio::log.info "Running upgrade commands..."
php /var/www/firefly/artisan firefly-iii:transaction-identifiers
php /var/www/firefly/artisan firefly-iii:migrate-to-groups
php /var/www/firefly/artisan firefly-iii:account-currencies
php /var/www/firefly/artisan firefly-iii:transfer-currencies
php /var/www/firefly/artisan firefly-iii:other-currencies
php /var/www/firefly/artisan firefly-iii:migrate-notes
php /var/www/firefly/artisan firefly-iii:migrate-attachments
php /var/www/firefly/artisan firefly-iii:bills-to-rules
php /var/www/firefly/artisan firefly-iii:bl-currency
php /var/www/firefly/artisan firefly-iii:cc-liabilities
php /var/www/firefly/artisan firefly-iii:back-to-journals
php /var/www/firefly/artisan firefly-iii:rename-account-meta
php /var/www/firefly/artisan firefly-iii:migrate-recurrence-meta

bashio::log.info "Running verification commands..."
php /var/www/firefly/artisan firefly-iii:fix-piggies
php /var/www/firefly/artisan firefly-iii:create-link-types
php /var/www/firefly/artisan firefly-iii:create-access-tokens
php /var/www/firefly/artisan firefly-iii:remove-bills
php /var/www/firefly/artisan firefly-iii:enable-currencies
php /var/www/firefly/artisan firefly-iii:fix-transfer-budgets
php /var/www/firefly/artisan firefly-iii:fix-uneven-amount
php /var/www/firefly/artisan firefly-iii:delete-zero-amount
php /var/www/firefly/artisan firefly-iii:delete-orphaned-transactions
php /var/www/firefly/artisan firefly-iii:delete-empty-journals
php /var/www/firefly/artisan firefly-iii:delete-empty-groups
php /var/www/firefly/artisan firefly-iii:fix-account-types
php /var/www/firefly/artisan firefly-iii:rename-meta-fields
php /var/www/firefly/artisan firefly-iii:fix-ob-currencies
php /var/www/firefly/artisan firefly-iii:fix-long-descriptions
php /var/www/firefly/artisan firefly-iii:fix-recurring-transactions

bashio::log.info "Running report commands..."
php /var/www/firefly/artisan firefly-iii:report-empty-objects
php /var/www/firefly/artisan firefly-iii:report-sum

php /var/www/firefly/artisan firefly-iii:restore-oauth-keys
php /var/www/firefly/artisan firefly-iii:set-latest-version --james-is-cool
php /var/www/firefly/artisan cache:clear > /dev/null 2>&1
php /var/www/firefly/artisan config:cache > /dev/null 2>&1

bashio::log.info "Cleanup cache and logs"
rm -rf /var/www/firefly/storage/framework/cache/data/*
rm -f /var/www/firefly/storage/logs/*.log

bashio::log.info "Starting PHP-FPM..."

exec php-fpm8.4 --nodaemonize
